# deployment/student-tracker.service
[Unit]
Description=STCM140 Student Tracker Dashboard
After=network.target

[Service]
Type=simple
User=jamditis
WorkingDirectory=/home/jamditis/projects/class
Environment="PATH=/home/jamditis/projects/class/venv/bin"
EnvironmentFile=/home/jamditis/.claude/.env
ExecStart=/home/jamditis/projects/class/venv/bin/python -m student_tracker.cli dashboard --port 5002
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=student-tracker

[Install]
WantedBy=multi-user.target
